# ansible/roles/r730xd-openbao/templates/openbao-auto-unseal.service.j2
# Managed by Ansible — r730xd-openbao role
# Runs the scripted unseal after OpenBao is up. Oneshot + RemainAfterExit
# so `systemctl is-active` reports "active (exited)" between boots.
#
# Boot-time ordering:
#   network-online.target → docker.service → foundation-openbao.service
#     → openbao-wait-ready.sh (ExecStartPre) → openbao-auto-unseal.sh

[Unit]
Description=OpenBao scripted unseal (Infisical-backed)
Requires={{ openbao_server_unit_name }}
After={{ openbao_server_unit_name }} network-online.target
Wants=network-online.target
ConditionPathExists={{ openbao_infisical_auth_env_file }}
ConditionPathExists={{ openbao_unseal_script }}

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre={{ openbao_wait_ready_script }}
ExecStart={{ openbao_unseal_script }}
# The unseal script self-retries a bounded number of times to ride out a
# transient boot-time DNS/network race reaching Infisical (RETRY_MAX×RETRY_DELAY
# plus per-attempt network timeouts), so give it headroom before systemd kills
# it mid-retry. It still exits non-zero once the budget is exhausted.
TimeoutStartSec=300
# Restart=no is deliberate: after the in-script retry budget is spent we fail
# loudly so Prometheus alerts fire, rather than masking a real outage with an
# endless restart loop. Operator reruns `systemctl restart openbao-auto-unseal`
# after fixing the root cause.
Restart=no

[Install]
WantedBy=multi-user.target
